# services/registry/deploy/xs-registry-backup.service
# SPDX-FileCopyrightText: 2026 Leitwolf <xs-lang.chess031@slmails.com>
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=Back up the X# package registry database
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=postgres
Group=postgres
UMask=0077
ExecStart=/usr/bin/pg_dump --format=custom --compress=6 --file=/var/lib/pgsql/backups/xs-registry/xs_registry.dump.tmp xs_registry
ExecStartPost=/usr/bin/mv --force /var/lib/pgsql/backups/xs-registry/xs_registry.dump.tmp /var/lib/pgsql/backups/xs-registry/xs_registry.dump
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/pgsql/backups/xs-registry

