# services/registry/deploy/xs-registry-backup.timer
# SPDX-FileCopyrightText: 2026 Leitwolf <xs-lang.chess031@slmails.com>
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=Daily X# package registry database backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m
Unit=xs-registry-backup.service

[Install]
WantedBy=timers.target
